# deploy/robot_status.timer
[Unit]
Description=Run Robot Status Fetcher Every 20 Seconds

[Timer]
# Initial delay after boot
OnBootSec=1min
# Time between consecutive runs
OnUnitActiveSec=20s
# Ensures the timer stays accurate
AccuracySec=1s

[Install]
WantedBy=timers.target

